# ops/systemd/convert_america.timer
[Timer]
OnCalendar=*-*-* 17:15:00
Persistent=true
