# deploy/systemd/daily_analysis@.service
[Unit]
Description=Sample daily analysis service (%i)

[Service]
WorkingDirectory=/root/telegram-crypto-bot-github
ExecStart=/root/telegram-crypto-bot-github/.venv/bin/python -u /root/telegram-crypto-bot-github/daily_analysis.py %i
