# deploy/systemd/dev3-analysis-hourly@.service
[Unit]
Description=Dev3 spot analysis (region %i)

[Service]
WorkingDirectory=/root/telegram-crypto-bot-github
ExecStart=/root/telegram-crypto-bot-github/.venv/bin/python -u /root/telegram-crypto-bot-github/daily_analysis.py %i
