# deploy/systemd/dev3-convert@.service
[Unit]
Description=Dev3 convert cycle (region %i)

[Service]
WorkingDirectory=/root/telegram-crypto-bot-github
ExecStart=/root/telegram-crypto-bot-github/.venv/bin/python -u /root/telegram-crypto-bot-github/run_convert_trade.py %i
