# ops/systemd/dev3-convert@.service.d/70-cli-args.conf
[Service]
ExecStart=
ExecStart=/root/telegram-crypto-bot-github/.venv/bin/python -u /srv/dev3/bin/convert_cycle.py %i
