# systemd/crypto-bot.service
[Unit]
Description=Telegram GPT Crypto Bot
After=network.target

[Service]
WorkingDirectory=/root/telegram-crypto-bot-github
EnvironmentFile=/root/telegram-crypto-bot-github/.env
ExecStart=/usr/bin/python3 main.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
