# systemd/crypto-bot.service
[Unit]
Description=Telegram Crypto Bot
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/crypto-bot
ExecStart=/usr/bin/python3 /opt/crypto-bot/main.py
Restart=no
ExecStopPost=/bin/sleep 1

[Install]
WantedBy=multi-user.target
