# systemd/crypto-bot.service
[Unit]
Description=Telegram GPT Crypto Bot
After=network.target

[Service]
Type=simple
WorkingDirectory=/root/telegram-crypto-bot-github
ExecStart=/usr/bin/python3 main.py
Restart=always
EnvironmentFile=/root/telegram-crypto-bot-github/systemd/crypto-bot.env
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
