# systemd/gpt_convert_asia.service
[Unit]
Description=Daily analysis for convert (asia session)
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/crypto-bot
ExecStart=/usr/bin/python3 /opt/crypto-bot/daily_analysis.py --mode convert --session asia
Restart=no
ExecStopPost=/bin/sleep 1

[Install]
WantedBy=multi-user.target
